# profile.d/50_android.sh
#!/bin/sh -e

prepend_path ~/repo/android-studio/bin
