# profile.d/50_android.sh
#!/bin/sh -e

PATH=~/repo/android-studio/bin:$PATH
export PATH
