# profile.d/50_apitrace.sh
#!/bin/sh -e

base=~/repo/apitrace

prepend_path $base/build
prepend_path $base/scripts
