# profile.d/50_apitrace.sh
#!/bin/sh -e

prepend_path ~/repo/apitrace/install/bin
