# profile.d/50_clion.sh
#!/bin/sh -e

prepend_path ~/repo/clion/bin
