# profile.d/50_clion.sh
#!/bin/sh -e

PATH=~/repo/clion/bin:$PATH
export PATH
