# profile.d/50_depot_tools.sh
#!/bin/sh -e

prepend_path ~/repo/depot_tools
