# profile.d/50_grep.sh
#!/bin/sh -e

if [ $( uname -s ) = Darwin ]; then
    alias grep='grep --color=auto'
fi
