# profile.d/50_go.sh
#!/bin/sh -e

prepend_path ~/repo/go/bin
