# profile.d/50_glsl.sh
#!/bin/sh -e

prepend_path ~/repo/glslang/build/StandAlone
prepend_path ~/repo/glslViewer
