# profile.d/50_gmock.sh
#!/bin/sh -e

PATH=$PATH:~/repo/gmock
export PATH
