# profile.d/50_j.sh
#!/bin/sh -e

J2_DIR=~/repo/j2
export JPY=$J2_DIR/j.py
. $J2_DIR/j.sh
