# profile.d/50_lefthook.sh
#!/bin/sh -e

prepend_path ~/repo/lefthook
