# profile.d/50_peat.sh
#!/bin/sh -e

PATH=~/repo/peat:$PATH
export PATH
