# profile.d/50_t.sh
#!/bin/sh -e

t()
{
    python ~/repo/t/t.py --task-dir ~ --list TODO $@
}
