# profile.d/50_tup.sh
#!/bin/sh -e

PATH=~/repo/tup:$PATH
export PATH
