# profile.d/50_tmuxp.sh
#!/bin/sh -e

hash tmuxp.bash 2>/dev/null && . tmuxp.bash
