# profile.d/50_z.sh
#!/bin/sh -e

if [ -f ~/repo/z/z.sh ]; then
    . ~/repo/z/z.sh
fi
