# profile.d/60_bin.sh
#!/bin/sh -e

PATH=~/bin:$PATH
export PATH
