# profile.d/90_lib.sh
#!/bin/sh -e

unset append_path
unset prepend_path
